# services/nomad-telegram-bot/ops/nomad-telegram-bot.service
[Unit]
Description=Nomad Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/nomad/services/nomad-telegram-bot
Environment=NODE_ENV=production
EnvironmentFile=/etc/nomad/nomad-telegram-bot.env
ExecStart=/usr/bin/node /opt/nomad/services/nomad-telegram-bot/dist/index.js
Restart=always
RestartSec=5
User=nomad
Group=nomad

[Install]
WantedBy=multi-user.target
